# roles/wireguard/files/etc/systemd/system/wireguard-config.timer
[Unit]
Description=Run WireGuard config on Schedule

[Timer]
OnBootSec=15s
OnUnitActiveSec=3m

[Install]
WantedBy=timers.target
