# roles/unifi/files/etc/systemd/system/controld-workaround.timer
[Unit]
Description=Workaround for Controld bug

[Timer]
OnBootSec=15s
OnUnitActiveSec=3m

[Install]
WantedBy=timers.target
