# roles/kopia_server/files/systemd/system/kopia-maintenance.timer
[Unit]
Description=Run Kopia maintenance
Requires=kopia-maintenance.service

[Timer]
Unit=kopia-maintenance.service
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
