# roles/kopia_server/files/systemd/system/kopia-sync.path
[Unit]
Description=Run Kopia Sync on trigger
Requires=kopia-sync.service

[Path]
# %t is /run
PathChanged=%t/kopia-sync.trigger
Unit=kopia-sync.service
TriggerLimitIntervalSec=6h
TriggerLimitBurst=1
