# roles/kopia_server/files/systemd/system/kopia-exporter.path
[Unit]
Description=Run Kopia Exporter on trigger
Requires=kopia-exporter.service

[Path]
# %t is /run
PathChanged=%t/kopia-exporter.trigger
Unit=kopia-exporter.service
