# roles/syncthing/files/etc/systemd/system/syncthing@.service.d/override.conf
[Service]
AmbientCapabilities=CAP_CHOWN CAP_FOWNER
